# templates/gitserver/omegaup-gitserver.service.erb
[Unit]
Description=omegaUp git service
After=network.target<% if @local_database %> mysql.service nginx.service php7.4-fpm.service<% end %>
<% if @local_database %>Requires=mysql.service nginx.service php7.4-fpm.service<% end %>

[Service]
Type=notify
User=omegaup
Group=omegaup
ExecStart=/usr/bin/omegaup-gitserver
ExecReload=/bin/kill -s HUP $MAINPID
WorkingDirectory=/var/lib/omegaup
Restart=always

[Install]
WantedBy=multi-user.target
