# templates/grader/omegaup-grader.service.erb
[Unit]
Description=omegaUp grading service
<% if @local_database %>
Requires=mysql.service omegaup-gitserver.service
After=network.target mysql.service omegaup-gitserver.service
<% else %>
Requires=mysql.service omegaup-gitserver.service
After=omegaup-gitserver.service
<% end %>

[Service]
Type=notify
User=omegaup
Group=omegaup
ExecStart=/usr/bin/omegaup-grader
ExecReload=/bin/kill -s HUP $MAINPID
WorkingDirectory=/var/lib/omegaup
Restart=always

[Install]
WantedBy=multi-user.target
